# scripts/eturnal.openrc
#!/sbin/openrc-run

name="${RC_SVCNAME}"
command="{{eturnal_bin_prefix}}/bin/${name}ctl"
extra_started_commands="reload"
capabilities="^cap_net_bind_service"

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting $name"
	$command daemon
	eend $?
}

stop() {
	ebegin "Stopping $name"
	$command stop
	eend $?
}

reload() {
	ebegin "Reloading $name"
	$command reload
	eend $?
}
