# scripts/eturnal.openrc
#!/sbin/openrc-run

command="{{eturnal_bin_prefix}}/bin/${RC_SVCNAME}ctl"
extra_started_commands="reload"
capabilities="^cap_net_bind_service"

depend() {
	need net
	after firewall
}

start() {
	ebegin "Starting $RC_SVCNAME"
	$command daemon
	eend $?
}

stop() {
	ebegin "Stopping $RC_SVCNAME"
	$command stop
	eend $?
}

reload() {
	ebegin "Reloading $RC_SVCNAME"
	$command reload
	eend $?
}
